# ops/systemd/buscrawl-raw-parquet-compaction-catchup.service
[Unit]
Description=Buscrawl recoverable raw parquet compaction catchup
After=network-online.target buscrawl-stream.service
Wants=network-online.target

[Service]
Type=simple
User=eljah
WorkingDirectory=/home/eljah/apps/buscrawl
Environment=BUS_RAW_PARQUET_CATCHUP_MAX_BATCHES=0
Environment=BUS_RAW_PARQUET_CATCHUP_CLEANUP_EVERY_BATCHES=1
Environment=BUS_RAW_PARQUET_CATCHUP_SLEEP_SECONDS=15
Environment=BUS_RAW_PARQUET_CATCHUP_EMPTY_SLEEP_SECONDS=300
Environment=BUS_COMPACTED_PARQUET_BOOTSTRAP_MODE=full-history
Environment=BUS_COMPACTED_PARQUET_NEWEST_FIRST=false
Environment=BUS_COMPACTED_PARQUET_MAX_FILES_PER_RUN=20000
Environment=BUS_COMPACTED_PARQUET_OUTPUT_PARTITIONS=2
Environment=BUS_COMPACTED_PARQUET_MIN_BATCH_FILES=100
Environment=BUS_COMPACTED_PARQUET_MIN_BATCH_BYTES=67108864
Environment=BUS_COMPACTED_PARQUET_MAX_OPEN_LAG_MINUTES=60
Environment=BUS_COMPACTED_PARQUET_ALLOW_SMALL_TAIL=false
Environment=BUS_COMPACTED_PARQUET_THREADS=1
Environment=BUS_COMPACTED_PARQUET_STAGED_READ_CHUNK_FILES=50
Environment=BUS_COMPACTED_PARQUET_STAGED_READ_CHUNK_SLEEP_MILLIS=500
Environment=BUS_COMPACTED_PARQUET_LEGACY_MAX_FILES_PER_RUN=2000
Environment=BUS_COMPACTED_PARQUET_LEGACY_MIN_BATCH_FILES=1000
Environment=BUS_COMPACTED_PARQUET_LEGACY_MIN_BATCH_BYTES=67108864
Environment=BUS_COMPACTED_PARQUET_LEGACY_MAX_OPEN_LAG_MINUTES=60
Environment=BUS_COMPACTED_PARQUET_LEGACY_ALLOW_SMALL_TAIL=false
Environment=BUS_RAW_PARQUET_CLEANUP_DRY_RUN=false
Environment=BUS_RAW_PARQUET_CLEANUP_MAX_FILES=10000
Environment=BUS_RAW_PARQUET_CLEANUP_MAX_RECORDS=20000
Environment=BUS_RAW_PARQUET_CLEANUP_MIN_AGE_HOURS=1
ExecStart=/home/eljah/apps/buscrawl/bin/run-raw-parquet-compaction-catchup.sh
Restart=always
RestartSec=60s
Nice=19
IOSchedulingClass=idle
IOSchedulingPriority=7
TimeoutStopSec=300s
NoNewPrivileges=true
StandardOutput=append:/home/eljah/apps/buscrawl/logs/raw-parquet-compaction-catchup.log
StandardError=append:/home/eljah/apps/buscrawl/logs/raw-parquet-compaction-catchup.log

[Install]
WantedBy=multi-user.target
