# ops/systemd/buscrawl-raw-parquet-cleanup.timer
[Unit]
Description=Run Buscrawl raw parquet cleanup periodically

[Timer]
OnBootSec=10min
OnUnitActiveSec=5min
AccuracySec=30s
Persistent=true

[Install]
WantedBy=timers.target
