# ops/systemd/buscrawl-stream.service.d/anti-storm.conf
[Unit]
StartLimitIntervalSec=10min
StartLimitBurst=5

[Service]
RestartSec=30s
